# clamav-scan/media_scan/media_scan_service/media_scan.service
[Unit]
Description=media-change-service

[Service]
ExecStart=/opt/auto-clamIPS/auto-clamav/media_scan/media-li.sh