# clamav-scan/clamscan-root-month-week-Day-Hour-timer/clamscan-root.service
[Unit]
Description=clamscan2 service

[Service]
ExecStart=/opt/auto-clamIPS/auto-clamav/clamav-root-scan.sh
