# clamav-scan/clamscan-root-month-week-Day-Hour-timer/clamscan-root-week.timer
[Unit]
Description=clamscan2 service timer

[Timer]
OnCalendar=Sat 3:00:00
Unit=clamscan-root.service

[Install]
WantedBy=multi-user.target

