# maltrail/maltrail-services/flush_blacklists.timer
[Unit]
Description=flush_blacklists timer

[Timer]
OnCalendar=Sat 3:00:00
Unit=flush_blacklists.service

[Install]
WantedBy=multi-user.target

