# maltrail/maltrail-services/flush_blacklists.service
[Unit]
Description=flush_blacklists service

[Service]
ExecStart=/opt/auto-clamIPS/maltrail/flush_blacklists.sh
