# maltrail/maltrail-services/listener_maltrail.service
[Unit]
Description=listener_maltrail service

[Service]
ExecStart=/opt/auto-clamIPS/maltrail/listener_maltrail.sh
