# maltrail/maltrail-services/maltrail_scan.service
[Unit]
Description=maltrail_scan service

[Service]
ExecStart=python3 /opt/auto-clamIPS/maltrail/maltrail_scan.py
