# clamav-scan/change_service/if-change-scan.service
[Unit]
Description=if-change-scan service

[Service]
ExecStart=python3 /opt/auto-clamIPS/auto-clamav/if-change-scan.py
