# clamav-scan/change_service/if-change.service
[Unit]
Description=if-change service

[Service]
ExecStart=/opt/auto-clamIPS/auto-clamav/if-change.sh
