# notify-send/notify-services/notify-send.timer
[Unit]
Description=notify-send service timer

[Timer]
OnCalendar=*:*:0/01
Unit=notify-send.service

[Install]
WantedBy=multi-user.target

