# rkhunter/rkhunter_service/rkhunter_scanner.service
[Unit]
Description=rkhunter_scanner service

[Service]
ExecStart=/opt/auto-clamIPS/rkhunter/rkhunter_scanner.sh
